# docker/cron.daily/daily
#!/bin/bash -l

# source /.env && cd /usr/src/client && yarn e2e
echo "hi"
